# files-r4s/etc/uci-defaults/99-custom-defaults
#!/bin/sh

uci set network.lan.ipaddr='192.168.50.254/24'
uci commit network

sed -i 's#^root:[^:]*:#root::#' /etc/shadow

exit 0
