# files-x86_64/etc/uci-defaults/99-custom-defaults
#!/bin/sh

uci set network.lan.ipaddr='10.10.10.123/24'
uci commit network

sed -i 's#^root:[^:]*:#root::#' /etc/shadow

exit 0
